# sysext-creator-deploy.path
[Unit]
Description=Monitor sysext staging directory

[Path]
PathChanged=/var/tmp/sysext-staging
MakeDirectory=yes
# FIX: Added Sticky Bit (1) to prevent symlink attacks and LPE vulnerabilities
DirectoryMode=1777

[Install]
WantedBy=multi-user.target
