# sysext-creator.service
[Unit]
Description=Sysext Creator Varlink Daemon
Requires=sysext-creator.socket
After=network.target

[Service]
ExecStart=/usr/bin/python3 /usr/libexec/sysext-creator/sysext-creator-daemon.py
StandardInput=socket
StandardOutput=journal

# User and Group
User=root
Group=root

# File System Isolation
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/extensions /var/lib/confexts /run/sysext-creator
ReadOnlyPaths=/etc /usr

# Privileges and Sandbox
NoNewPrivileges=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes

# Network and Syscalls
RestrictAddressFamilies=AF_UNIX
SystemCallArchitectures=native
SystemCallFilter=@system-service @file-system @process

[Install]
WantedBy=multi-user.target
